# experiments/systemd_units/stb_game_ranking_pusher.service
[Unit]
Description=Super Tilt Bro. game ranking pusher
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=60
User=debian
ExecStartPre=
ExecStart=PATH_TO/game_ranking_pusher --ranking-server-addr 127.0.0.1 --replay-server-addr 127.0.0.1 --bmov-dir /usr/share/stb/games --game-server-name test-server
ExecStartPost=
ExecStop=
ExecReload=

[Install]
WantedBy=multi-user.target
