# experiments/systemd_units/webrtc_to_udp.service
[Unit]
Description=WebRTC and UDP cross communication server
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=10
User=debian
ExecStartPre=
ExecStart=/usr/bin/node PATH_TO/server.js
ExecStartPost=
ExecStop=
ExecReload=

[Install]
WantedBy=multi-user.target
